# spot.service
[Unit]
Description=Telegram Bot
After=network.target

[Service]
User=user_name
Group=user_name

WorkingDirectory=/home/user_name/spot_signals/
Environment="PYTHONPATH=/home/user_name/spot_signals/"
ExecStart=/home/user_name/spot_signals/venv/bin/python /home/user_name/spot_signals/scan.py

[Install]
WantedBy=multi-user.target
